# root/etc/services.d/openssh/run
#!/usr/bin/with-contenv bash

exec /usr/sbin/sshd -D -E /config/log/sshd/auth.log
